# lib/support/upstart/gitlab-ci-runner.conf
# Upstart configuration file for gitlab-ci-runner
# Control with:
#
#   sudo service gitlab-ci-runner start
#

description "GitLab CI Runner"
author "GitLab B.V."

start on runlevel [2345]

# Always restart the runner if it dies / exits
respawn

# Log to /var/log/upstart/
console log

# Run the process as gitlab_ci_runner:gitlab_ci_runner
setuid gitlab_ci_runner
setgid gitlab_ci_runner

# This is where you installed the gitlab-ci-runner source code
chdir /home/gitlab_ci_runner/gitlab-ci-runner

script
  # optional: set PATH or other environment variables
  # export PATH=/home/gitlab_ci_runner/.rvm/rubies/default/bin:$PATH
  export HOME=/home/gitlab_ci_runner
  exec bundle exec bin/runner
end script
